# tools/demo/deploy/yamlover-demo.service
# systemd unit for the demo server. Install to /etc/systemd/system/, then:
#   systemctl daemon-reload && systemctl enable --now yamlover-demo
# Adjust WorkingDirectory to wherever the repo is checked out on the VM.

[Unit]
Description=yamlover demo server
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/yamlover/tools/demo
ExecStart=/usr/bin/node bin/demo-server.js
Restart=always
RestartSec=3
# Secrets + config (DEMO_DRIVER, DEMO_BASE_URL, RESEND_API_KEY, …) live here.
EnvironmentFile=/etc/yamlover-demo.env
# Persistent spool + sqlite across reboots (set DB_PATH/SPOOL_DIR under it in the env file).
StateDirectory=yamlover-demo
ReadWritePaths=/var/lib/yamlover-demo
# Light hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
